# sysutils/py-zettarepl/files/zettarepl.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: zettarepl
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# zettarepl_enable (bool): Set to NO by default.
#                          Set it to YES to enable zettarepl.
# zettarepl_config (path): Set to %%PREFIX%%/etc/zettarepl/zettarepl.yaml
#                          by default.

. /etc/rc.subr

name=zettarepl
rcvar=zettarepl_enable

load_rc_config $name

: ${zettarepl_enable:="NO"}
: ${zettarepl_config="%%PREFIX%%/etc/zettarepl/zettarepl.yaml"}

command=%%PREFIX%%/bin/${name}
command_args="-f $zettarepl_config"
command_interpreter="%%PYTHON_CMD%%"

run_rc_command "$1"
